# vm_escape_through_nfc/etc/systemd/system/nfc_reader.service
[Unit]
Description=PN532 NFC Reader Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nfc_reader
Restart=on-failure

[Install]
WantedBy=multi-user.target
